# lily-bot-autoupdate.service
[Unit]
Description=Lily Music Bot Auto-Update & Self-Heal
Documentation=https://github.com/nishkarshk212/lily_music-

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/lily_music-
Environment="PATH=/root/lily_music-/venv/bin"
ExecStart=/bin/bash /root/lily_music-/auto-update.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lily-bot-autoupdate

[Install]
WantedBy=multi-user.target
